# cloxhack-supervisor.conf
[program:cloxhack]
command=/home/pi/.virtualenvs/cloxhack/bin/python /home/pi/cloxhack/clock.py
startsecs=10
user=pi
